# pg-backup.timer
[Unit]
Description=PostgreSQL Backup Timer - Daily at 12:00 PM UTC+3
Requires=pg-backup.service

[Timer]
# Run at 12:00 PM UTC+3 (09:00 UTC)
OnCalendar=*-*-* 09:00:00

# If system was powered off during scheduled time, run backup on next boot
Persistent=true

# Accuracy window (backup can run within 1 minute of scheduled time)
AccuracySec=1min

[Install]
WantedBy=timers.target
