# pg-backup.timer
[Unit]
Description=PostgreSQL Backup Timer - Daily at Midnight UTC+3
Requires=pg-backup.service

[Timer]
# Run at 00:00 (midnight) UTC+3 (21:00 UTC previous day)
OnCalendar=*-*-* 00:00:00

# If system was powered off during scheduled time, run backup on next boot
Persistent=true

# Accuracy window (backup can run within 1 minute of scheduled time)
AccuracySec=1min

[Install]
WantedBy=timers.target